# .root/fs/etc/systemd/user/wayland-session.target
[Unit]
Description=Current User Wayland Session
Requires=basic.target
RefuseManualStart=yes
# StopWhenUnneeded=yes
BindsTo=graphical-session.target
